# coordinator_service/coordinator.service
[Unit]
Description=Coordinator Service for Smart Mirror

[Service]
ExecStart=/usr/local/bin/coordinator

[Install]
WantedBy=multi-user.target
